# supervisor/component.conf
; supervisord program block for EthernetIpAdapter on the HOST platform.
;
; NOTE: in a rendered deployment this file is OWNED BY THE DEPLOYMENT RENDERER — the EdgeCommons
; Deployment Studio HOST renderer regenerates it from the deployment model (install paths, config
; paths, thing name, user, log dirs). Hand-edits are overwritten on the next render. Treat the copy
; in this repo as the template default / a standalone reference you can run yourself:
;
;   sudo cp supervisor/component.conf /etc/supervisor/conf.d/EthernetIpAdapter.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl status EthernetIpAdapter
;
; Expected layout (what the renderer lays down):
;   /opt/edgecommons/EthernetIpAdapter/bin/ethernetipadapter        the release binary (cargo build --release)
;   /etc/edgecommons/EthernetIpAdapter/config.json            the component config    (-c FILE)
;   /etc/edgecommons/EthernetIpAdapter/messaging.json         the MQTT broker config  (--transport MQTT)
;   /var/log/edgecommons/                                     stdout/stderr logs

[program:EthernetIpAdapter]
command=/opt/edgecommons/EthernetIpAdapter/bin/ethernetipadapter --platform HOST --transport MQTT /etc/edgecommons/EthernetIpAdapter/messaging.json -c FILE /etc/edgecommons/EthernetIpAdapter/config.json -t ethernetipadapter-dev
directory=/opt/edgecommons/EthernetIpAdapter
user=edgecommons
; The `logging` config block drives the component's own log level; RUST_BACKTRACE only affects panics.
environment=RUST_BACKTRACE="1"

autostart=true
autorestart=true
startsecs=5
startretries=3
priority=200

; The library installs a SIGTERM handler that runs its graceful, idempotent shutdown (flip /readyz
; to 503, unsubscribe, close messaging/metrics/streams). Give it room to finish.
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true

stdout_logfile=/var/log/edgecommons/EthernetIpAdapter.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/edgecommons/EthernetIpAdapter.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
redirect_stderr=false
